# deploy/systemd/official-sources-hermes-auditor.timer
[Unit]
Description=Run official-sources Hermes strict release auditor

[Timer]
OnCalendar=*-*-* 00:15:00 UTC
Persistent=true
Unit=official-sources-hermes-auditor.service

[Install]
WantedBy=timers.target
