# deploy/systemd/official-sources-hermes-auditor.service
[Unit]
Description=official-sources Hermes strict release auditor
Documentation=https://github.com/example/official-sources
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=official-sources
Group=official-sources
SupplementaryGroups=systemd-journal
WorkingDirectory=/opt/official-sources/app
EnvironmentFile=/opt/official-sources/.env
Environment=APP_ROOT=/opt/official-sources/app
Environment=STATE_ROOT=/var/lib/hermes-official-sources-auditor
Environment=RELEASE_CONTRACT=/etc/official-sources/hermes-audit-contract.yaml
Environment=OFFICIAL_SOURCES_BIN=/opt/official-sources/app/.venv/bin/official-sources
ExecStart=/opt/hermes-official-sources-auditor/bin/run-official-sources-hermes-auditor.sh
NoNewPrivileges=true
PrivateTmp=true
